# yhs_说明/yhs_fr07_pro_oc/can卡设置说明/opencan.rules
ACTION=="add",SUBSYSTEMS=="usb",ATTRS{idVendor}=="1d50",ATTRS{idProduct}=="606f",RUN+="/bin/bash /opt/mnt/opencan.sh"
